# _content/journal/2020/service_management_with_systemd/myservice.path
[Unit]
Description="Watch /data/app/bin/myservice"

[Path]
PathChanged=/data/app/bin/myservice
Unit=systemctl-restart@%p.service

[Install]
WantedBy=multi-user.target
