# _content/journal/2020/service_management_with_systemd/systemctl-restart@.service
[Unit]
Description=systemctl-restart@%i

[Service]
ExecStart=/bin/systemctl restart %i
