# cmd/kilabit/kilabit.path
[Unit]
Description="Watch kilabit"

[Path]
PathChanged=/home/ms/go/bin/kilabit
Unit=systemctl-restart@%p.service

[Install]
WantedBy=multi-user.target
