# cmd/kilabit/kilabit.service
[Unit]
Description=kilabit

[Service]
Type=simple
ExecStart=/home/ms/go/bin/kilabit

[Install]
WantedBy=default.target
